# deploy/truhlik.service
[Unit]
Description=Truhlik Controller
After=network.target

[Service]
User=pavelhalbich
WorkingDirectory=/opt/truhlik/app
ExecStart=/opt/truhlik/venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
